# cloud/assets/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:open-xchange]
command=/data/open-xchange.sh
