# cloud/assets/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:postfix]
command=/data/open-xchange.sh
